# DianaStack/roles/diana-monitor/templates/supervisor-dmon.conf
[program:d-mon]
command=python /usr/src/DIANA/utilities/connect/d-mon.py --config /tmp/dmon_cfg.yml
autostart=true
autorestart=true
redirect_stderr=true